# crontab
@reboot  /usr/bin/python3 /app/app.py >> /var/log/cron.log 2>&1
0 8 * * * /usr/bin/python3 /app/app.py >> /var/log/cron.log 2>&1
